# cloudflared-cric-photos.service
[Unit]
Description=Cloudflare Tunnel for CRIC Photos
After=network.target

[Service]
Type=simple
User=admin
ExecStart=/usr/local/bin/cloudflared tunnel --config /home/admin/.cloudflared/config.yml run
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
